# data/cloud-init.service.in
[Unit]
Description=cloud-init job
After=local-fs.target network.target
Before=sshd.service sshd-keygen.service
Requires=network.target
Wants=local-fs.target sshd.service sshd-keygen.service

[Service]
Type=oneshot
ExecStart=@prefix@/bin/cloud-init --fix-disk --metadata --user-data-once
RemainAfterExit=yes
TimeoutSec=0
PrivateTmp=yes

# Output needs to appear in instance console output
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
