# data/cloud-init-pre-network.service.in
[Unit]
Description=cloud-init pre-network job
After=local-fs.target
Before=cloud-init.service network.target systemd-networkd.service
Wants=local-fs.target sshd.service sshd-keygen.service

[Service]
Type=oneshot
ExecStart=@prefix@/bin/cloud-init --first-boot-setup --no-network --metadata
RemainAfterExit=yes
TimeoutSec=0
PrivateTmp=yes

# Output needs to appear in instance console output
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
