# data/ucd.service.in
[Unit]
Description=micro-config-drive job
After=multi-user.target
Before=network.target systemd-networkd.service
Wants=local-fs.target sshd.service sshd-keygen.service

[Service]
Type=oneshot
ExecStart=@prefix@/bin/ucd --fix-disk --first-boot-setup --metadata --user-data-once
RemainAfterExit=yes
TimeoutSec=0

# Output needs to appear in instance console output
StandardOutput=journal+console

[Install]
